# docs/notes-clipboard-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.notes-clipboard-sync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd /Users/mimiko/Projects/notes-clipboard-sync &amp;&amp; pnpm notes once</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/notes-clipboard-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/notes-clipboard-sync.err.log</string>
  </dict>
</plist>
